# packing/deb/prerm
#!/bin/bash

systemctl stop webhooks

systemctl disable webhooks

exit 0